# modules/snmpd/files/etc/uci-defaults/99-snmpd
# uci batch << EOF
# commit
# EOF
# fi
uci batch << EOF

set snmpd.@agent[0]=agent
set snmpd.@agent[0].agentaddress='UDP:161,UDP6:161'
set snmpd.@agentx[0]=agentx
set snmpd.@agentx[0].agentxsocket='/var/run/agentx.sock'
set snmpd.public=com2sec
set snmpd.public.secname='ro'
set snmpd.public.source='default'
set snmpd.public.community='public'
set snmpd.private=com2sec
set snmpd.private.secname='rw'
set snmpd.private.source='localhost'
set snmpd.private.community='private'
set snmpd.public6=com2sec6
set snmpd.public6.secname='ro'
set snmpd.public6.source='default'
set snmpd.public6.community='public'
set snmpd.private6=com2sec6
set snmpd.private6.secname='rw'
set snmpd.private6.source='localhost'
set snmpd.private6.community='private'
set snmpd.public_v1=group
set snmpd.public_v1.group='public'
set snmpd.public_v1.version='v1'
set snmpd.public_v1.secname='ro'
set snmpd.public_v2c=group
set snmpd.public_v2c.group='public'
set snmpd.public_v2c.version='v2c'
set snmpd.public_v2c.secname='ro'
set snmpd.public_usm=group
set snmpd.public_usm.group='public'
set snmpd.public_usm.version='usm'
set snmpd.public_usm.secname='ro'
set snmpd.private_v1=group
set snmpd.private_v1.group='private'
set snmpd.private_v1.version='v1'
set snmpd.private_v1.secname='rw'
set snmpd.private_v2c=group
set snmpd.private_v2c.group='private'
set snmpd.private_v2c.version='v2c'
set snmpd.private_v2c.secname='rw'
set snmpd.private_usm=group
set snmpd.private_usm.group='private'
set snmpd.private_usm.version='usm'
set snmpd.private_usm.secname='rw'
set snmpd.all=view
set snmpd.all.viewname='all'
set snmpd.all.type='included'
set snmpd.all.oid='.1'
set snmpd.public_access=access
set snmpd.public_access.group='public'
set snmpd.public_access.context='none'
set snmpd.public_access.version='any'
set snmpd.public_access.level='noauth'
set snmpd.public_access.prefix='exact'
set snmpd.public_access.read='all'
set snmpd.public_access.write='none'
set snmpd.public_access.notify='none'
set snmpd.private_access=access
set snmpd.private_access.group='private'
set snmpd.private_access.context='none'
set snmpd.private_access.version='any'
set snmpd.private_access.level='noauth'
set snmpd.private_access.prefix='exact'
set snmpd.private_access.read='all'
set snmpd.private_access.write='all'
set snmpd.private_access.notify='all'
set snmpd.@system[0]=system
set snmpd.@system[0].sysLocation='office'
set snmpd.@system[0].sysContact='bofh@example.com'
set snmpd.@system[0].sysName='HeartOfGold'
set snmpd.@exec[0]=exec
set snmpd.@exec[0].name='filedescriptors'
set snmpd.@exec[0].prog='/bin/cat'
set snmpd.@exec[0].args='/proc/sys/fs/file-nr'
set snmpd.@engineid[0]=engineid
set snmpd.@engineid[0].engineidtype='3'
set snmpd.@engineid[0].engineidnic='eth0'
set snmpd.general=snmpd
set snmpd.general.enabled='1'

commit
EOF
sed -i '/^service_triggers(){/,/^}/ s/^/# /' /etc/init.d/snmpd