# modules/prefer-ipv6-settings/files/etc/uci-defaults/99-ipv6
uci batch << EOF

delete network.globals.ula_prefix
set network.globals.ula_prefix='fdd6:4872:af1f::/48'

commit
EOF