# modules/openclash/files/etc/uci-defaults/90-openclash
while uci -q delete openclash.@dns_servers[0]; do :; done
while uci -q delete openclash.@config_subscribe[0]; do :; done

# 默认 SOCKS 不需要认证
while uci -q delete openclash.@authentication[0]; do :; done

username=$AUTH_USERNAME
password=$AUTH_PASSWORD
subscribe_addr=$SUBSCRIBE

if [ -z "$username" ]; then
    username='brainiac'
fi
if [ -z "$password" ]; then
    password='123456'
fi
if [ -z "$subscribe_addr" ]; then
    subscribe_addr='https://example.com/subscribe'
fi
 
uci batch << EOF

set openclash.config=openclash
set openclash.config.proxy_port='7892'
set openclash.config.tproxy_port='7895'
set openclash.config.mixed_port='7893'
set openclash.config.socks_port='7891'
set openclash.config.http_port='7890'
set openclash.config.dns_port='7874'
set openclash.config.update='0'
set openclash.config.cn_port='9090'
set openclash.config.dashboard_password='123456'
set openclash.config.dashboard_forward_ssl='0'
set openclash.config.rule_source='0'
set openclash.config.enable_custom_clash_rules='0'
set openclash.config.dns_advanced_setting='0'
set openclash.config.servers_if_update='0'
set openclash.config.servers_update='0'
set openclash.config.log_level='silent'
set openclash.config.proxy_mode='rule'
set openclash.config.intranet_allowed='1'
set openclash.config.disable_udp_quic='1'
set openclash.config.enable_rule_proxy='0'
set openclash.config.filter_aaaa_dns='0'
set openclash.config.small_flash_memory='0'
set openclash.config.interface_name='0'
set openclash.config.log_size='1024'
set openclash.config.store_fakeip='1'
set openclash.config.custom_fallback_filter='0'
set openclash.config.append_wan_dns='1'
set openclash.config.stream_domains_prefetch='0'
set openclash.config.stream_auto_select='0'
set openclash.config.dns_remote='1'
set openclash.config.github_address_mod='0'
set openclash.config.router_self_proxy='1'
set openclash.config.release_branch='master'
set openclash.config.chnr6_custom_url='https://ispip.clang.cn/all_cn_ipv6.txt'
set openclash.config.default_resolvfile='/tmp/resolv.conf.d/resolv.conf.auto'
set openclash.config.auto_update='1'
set openclash.config.config_path='/etc/openclash/config/DUANG.yaml'
set openclash.config.restricted_mode='0'
set openclash.config.create_config='0'
set openclash.config.chnr_auto_update='1'
set openclash.config.chnr_update_week_time='1'
set openclash.config.auto_restart='0'
set openclash.config.auto_restart_week_time='1'
set openclash.config.auto_restart_day_time='0'
set openclash.config.enable_custom_dns='1'
set openclash.config.chnr_update_day_time='4'
set openclash.config.dnsmasq_resolvfile='/tmp/resolv.conf.auto'
set openclash.config.dnsmasq_noresolv='0'
set openclash.config.enable_redirect_dns='0'
set openclash.config.redirect_dns='0'
set openclash.config.operation_mode='redir-host'
set openclash.config.common_ports='0'
set openclash.config.append_default_dns='1'
set openclash.config.chnr_custom_url='https://ispip.clang.cn/all_cn_cidr.txt'
set openclash.config.core_version='linux-amd64'
set openclash.config.cndomain_custom_url='https://testingcf.jsdelivr.net/gh/felixonmars/dnsmasq-china-list@master/accelerated-domains.china.conf'
set openclash.config.ipv6_dns='1'
set openclash.config.ipv6_enable='1'
set openclash.config.china_ip6_route='1'
set openclash.config.china_ip_route='1'
set openclash.config.bypass_gateway_compatible='0'
set openclash.config.config_auto_update_mode='0'
set openclash.config.cachesize_dns='0'
set openclash.config.other_rule_auto_update='0'
set openclash.config.en_mode='redir-host'
set openclash.config.enable_udp_proxy='1'
set openclash.config.enable_meta_core='1'
set openclash.config.enable_tcp_concurrent='1'
set openclash.config.find_process_mode='off'
set openclash.config.client_fingerprint='0'
set openclash.config.enable_meta_sniffer='1'
set openclash.config.enable_meta_sniffer_pure_ip='1'
set openclash.config.enable_meta_sniffer_custom='0'
set openclash.config.geodata_loader='standard'
set openclash.config.enable_geoip_dat='1'
set openclash.config.geoip_auto_update='1'
set openclash.config.geoip_update_day_time='4'
set openclash.config.geoip_custom_url='https://testingcf.jsdelivr.net/gh/Loyalsoldier/v2ray-rules-dat@release/geoip.dat'
set openclash.config.core_type='Meta'
set openclash.config.geosite_auto_update='1'
set openclash.config.geosite_update_day_time='4'
set openclash.config.geosite_custom_url='https://testingcf.jsdelivr.net/gh/Loyalsoldier/v2ray-rules-dat@release/geosite.dat'
set openclash.config.disable_masq_cache='0'
set openclash.config.enable_custom_domain_dns_server='0'
set openclash.config.lan_ac_mode='0'
set openclash.config.geo_auto_update='1'
set openclash.config.geo_update_week_time='1'
set openclash.config.geo_update_day_time='4'
set openclash.config.geo_custom_url='http://www.ideame.top/mmdb/Country.mmdb'
set openclash.config.geoip_update_week_time='1'
set openclash.config.geosite_update_week_time='1'
set openclash.config.config_update_week_time='1'
set openclash.config.auto_update_time='3'
set openclash.config.tolerance='100'
set openclash.config.urltest_interval_mod='300'
set openclash.config.urltest_address_mod='http://cp.cloudflare.com/generate_204'
set openclash.config.custom_name_policy='0'
set openclash.config.custom_host='0'
set openclash.config.global_client_fingerprint='0'
set openclash.config.delay_start='30'
set openclash.config.dashboard_type='Meta'
set openclash.config.enable='1'
set openclash.config.yacd_type='Meta'
set openclash.config.restart='0'
set openclash.config.config_reload='1'
add openclash dns_servers
set openclash.@dns_servers[-1].group='nameserver'
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].enabled='1'
set openclash.@dns_servers[-1].ip='127.0.0.1'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
set openclash.@dns_servers[-1].port='6053'
add openclash dns_servers
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='127.0.0.1'
set openclash.@dns_servers[-1].port='6553'
set openclash.@dns_servers[-1].enabled='1'
add openclash dns_servers
set openclash.@dns_servers[-1].group='nameserver'
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].ip='119.28.28.28'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='nameserver'
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].ip='223.5.5.5'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].type='https'
set openclash.@dns_servers[-1].ip='doh.pub/dns-query'
set openclash.@dns_servers[-1].group='nameserver'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].type='https'
set openclash.@dns_servers[-1].ip='dns.alidns.com/dns-query'
set openclash.@dns_servers[-1].group='nameserver'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
set openclash.@dns_servers[-1].ip='127.0.0.1'
set openclash.@dns_servers[-1].port='7553'
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].enabled='0'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='dns.google'
set openclash.@dns_servers[-1].port='853'
set openclash.@dns_servers[-1].type='tls'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].type='https'
set openclash.@dns_servers[-1].ip='1.1.1.1/dns-query'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='1.1.1.1'
set openclash.@dns_servers[-1].port='853'
set openclash.@dns_servers[-1].type='tls'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='8.8.8.8'
set openclash.@dns_servers[-1].port='853'
set openclash.@dns_servers[-1].type='tls'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='2001:4860:4860::8888'
set openclash.@dns_servers[-1].port='53'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='2001:4860:4860::8844'
set openclash.@dns_servers[-1].port='53'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].type='udp'
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='2001:da8::666'
set openclash.@dns_servers[-1].port='53'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].type='https'
set openclash.@dns_servers[-1].ip='public.dns.iij.jp/dns-query'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].type='https'
set openclash.@dns_servers[-1].ip='jp.tiar.app/dns-query'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].type='https'
set openclash.@dns_servers[-1].ip='jp.tiarap.org/dns-query'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='jp.tiar.app'
set openclash.@dns_servers[-1].type='tls'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash dns_servers
set openclash.@dns_servers[-1].group='fallback'
set openclash.@dns_servers[-1].ip='dot.tiar.app'
set openclash.@dns_servers[-1].type='tls'
set openclash.@dns_servers[-1].enabled='0'
set openclash.@dns_servers[-1].node_resolve='0'
set openclash.@dns_servers[-1].interface='Disable'
set openclash.@dns_servers[-1].specific_group='Disable'
add openclash config_subscribe
set openclash.@config_subscribe[-1].enabled='1'
set openclash.@config_subscribe[-1].name='SubConfig'
set openclash.@config_subscribe[-1].emoji='false'
set openclash.@config_subscribe[-1].udp='true'
set openclash.@config_subscribe[-1].skip_cert_verify='false'
set openclash.@config_subscribe[-1].sort='false'
set openclash.@config_subscribe[-1].node_type='false'
set openclash.@config_subscribe[-1].rule_provider='false'
set openclash.@config_subscribe[-1].sub_convert='0'
add_list openclash.@config_subscribe[-1].keyword='香港'
add_list openclash.@config_subscribe[-1].keyword='Hong Kong'
add_list openclash.@config_subscribe[-1].keyword='HK'
add_list openclash.@config_subscribe[-1].keyword='新加坡'
add_list openclash.@config_subscribe[-1].keyword='Singapore'
add_list openclash.@config_subscribe[-1].keyword='SG'
set openclash.@config_subscribe[-1].address='$subscribe_addr'
add openclash authentication
set openclash.@authentication[-1].enabled='1'
set openclash.@authentication[-1].username='$username'
set openclash.@authentication[-1].password='$password'
set openclash.config.operation_mode='redir-host'
set openclash.config.en_mode='redir-host'
set openclash.config.enable_respect_rules='1'
set openclash.config.custom_fallback_filter='1'
set openclash.config.ipv6_enable='1'
set openclash.config.china_ip6_route='1'
set openclash.config.china_ip_route='1'
set openclash.config.geo_auto_update='1'
set openclash.config.geo_update_day_time='6'
set openclash.config.geo_update_week_time='*'
set openclash.config.geoip_auto_update='1'
set openclash.config.geoip_update_week_time='*'
set openclash.config.geoip_update_day_time='4'
set openclash.config.auto_restart_week_time='*'
set openclash.config.chnr_update_day_time='3'
set openclash.config.chnr_update_week_time='*'
set openclash.config.chnr_auto_update='1'
set openclash.config.geosite_auto_update='1'
set openclash.config.enable_v6_udp_proxy='1'
set openclash.config.geosite_update_week_time='*'
set openclash.config.geosite_update_day_time='0'
set openclash.config.geosite_custom_url='https://testingcf.jsdelivr.net/gh/Loyalsoldier/v2ray-rules-dat@release/geosite.dat'
set openclash.config.auto_restart='1'
set openclash.config.auto_restart_day_time='3'
set openclash.config.enable_dns='1'
set openclash.config.ipv6_dns='1'
set openclash.config.append_wan_dns='0'
set openclash.config.disable_udp_quic='0'
set openclash.config.auto_update='1'
set openclash.config.config_auto_update_mode='0'
set openclash.config.config_update_week_time='*'

set openclash.config.core_type='Meta'
set openclash.config.enable_meta_core='1'
set openclash.config.enable_meta_sniffer='0'
set openclash.config.common_ports='0'
set openclash.config.append_default_dns='0'
set openclash.config.enable_custom_dns='1'

set openclash.config.find_process_mode='off'
commit
EOF