# modules/openclash/files/etc/uci-defaults/90-openclash
while uci -q delete openclash.@dns_servers[0]; do :; done
while uci -q delete openclash.@config_subscribe[0]; do :; done

# 默认 SOCKS 不需要认证
while uci -q delete openclash.@authentication[0]; do :; done

username=$AUTH_USERNAME
password=$AUTH_PASSWORD
subscribe_addr=$SUBSCRIBE

if [ -z "$username" ]; then
    username='brainiac'
fi
if [ -z "$password" ]; then
    password='123456'
fi
 
uci batch << EOF

set openclash.config=openclash
set openclash.config.proxy_port='7892'
set openclash.config.tproxy_port='7895'
set openclash.config.mixed_port='7893'
set openclash.config.socks_port='7891'
set openclash.config.http_port='7890'
set openclash.config.dns_port='7874'
set openclash.config.update='0'
set openclash.config.cn_port='9090'
set openclash.config.dashboard_password='123456'
set openclash.config.dashboard_forward_ssl='0'
set openclash.config.rule_source='0'
set openclash.config.enable_custom_clash_rules='0'
set openclash.config.dns_advanced_setting='0'
set openclash.config.servers_if_update='0'
set openclash.config.servers_update='0'
set openclash.config.log_level='silent'
set openclash.config.proxy_mode='rule'
set openclash.config.intranet_allowed='1'
set openclash.config.disable_udp_quic='1'
set openclash.config.enable_rule_proxy='0'
set openclash.config.filter_aaaa_dns='0'
set openclash.config.small_flash_memory='0'
set openclash.config.interface_name='0'
set openclash.config.log_size='1024'
set openclash.config.store_fakeip='1'
set openclash.config.custom_fallback_filter='0'
set openclash.config.append_wan_dns='1'
set openclash.config.stream_domains_prefetch='0'
set openclash.config.stream_auto_select='0'
set openclash.config.dns_remote='1'
set openclash.config.github_address_mod='0'
set openclash.config.router_self_proxy='1'
set openclash.config.release_branch='master'
set openclash.config.chnr6_custom_url='https://ispip.clang.cn/all_cn_ipv6.txt'
set openclash.config.default_resolvfile='/tmp/resolv.conf.d/resolv.conf.auto'
set openclash.config.auto_update='1'
set openclash.config.config_path='/etc/openclash/config/DUANG.yaml'
set openclash.config.restricted_mode='0'
set openclash.config.create_config='0'
set openclash.config.chnr_auto_update='1'
set openclash.config.chnr_update_week_time='1'
set openclash.config.auto_restart='0'
set openclash.config.auto_restart_week_time='1'
set openclash.config.auto_restart_day_time='0'
set openclash.config.enable_custom_dns='1'
set openclash.config.chnr_update_day_time='4'
set openclash.config.dnsmasq_resolvfile='/tmp/resolv.conf.auto'
set openclash.config.dnsmasq_noresolv='0'
set openclash.config.enable_redirect_dns='0'
set openclash.config.redirect_dns='0'
set openclash.config.operation_mode='redir-host'
set openclash.config.common_ports='0'
set openclash.config.append_default_dns='1'
set openclash.config.chnr_custom_url='https://ispip.clang.cn/all_cn_cidr.txt'
set openclash.config.core_version='linux-amd64'
set openclash.config.cndomain_custom_url='https://testingcf.jsdelivr.net/gh/felixonmars/dnsmasq-china-list@master/accelerated-domains.china.conf'
set openclash.config.ipv6_dns='1'
set openclash.config.ipv6_enable='1'
set openclash.config.china_ip6_route='1'
set openclash.config.china_ip_route='1'
set openclash.config.bypass_gateway_compatible='0'
set openclash.config.config_auto_update_mode='0'
set openclash.config.cachesize_dns='0'
set openclash.config.other_rule_auto_update='0'
set openclash.config.en_mode='redir-host'
set openclash.config.enable_udp_proxy='1'
set openclash.config.enable_meta_core='1'
set openclash.config.enable_tcp_concurrent='1'
set openclash.config.find_process_mode='off'
set openclash.config.client_fingerprint='0'
set openclash.config.enable_meta_sniffer='1'
set openclash.config.enable_meta_sniffer_pure_ip='1'
set openclash.config.enable_meta_sniffer_custom='0'
set openclash.config.geodata_loader='standard'
set openclash.config.enable_geoip_dat='1'
set openclash.config.geoip_auto_update='1'
set openclash.config.geoip_update_day_time='4'
set openclash.config.geoip_custom_url='https://testingcf.jsdelivr.net/gh/Loyalsoldier/v2ray-rules-dat@release/geoip.dat'
set openclash.config.core_type='Meta'
set openclash.config.geosite_auto_update='1'
set openclash.config.geosite_update_day_time='4'
set openclash.config.geosite_custom_url='https://testingcf.jsdelivr.net/gh/Loyalsoldier/v2ray-rules-dat@release/geosite.dat'
set openclash.config.disable_masq_cache='0'
set openclash.config.enable_custom_domain_dns_server='0'
set openclash.config.lan_ac_mode='0'
set openclash.config.geo_auto_update='1'
set openclash.config.geo_update_week_time='1'
set openclash.config.geo_update_day_time='4'
set openclash.config.geo_custom_url='http://www.ideame.top/mmdb/Country.mmdb'
set openclash.config.geoip_update_week_time='1'
set openclash.config.geosite_update_week_time='1'
set openclash.config.config_update_week_time='1'
set openclash.config.auto_update_time='3'
set openclash.config.tolerance='100'
set openclash.config.urltest_interval_mod='300'
set openclash.config.urltest_address_mod='http://cp.cloudflare.com/generate_204'
set openclash.config.custom_name_policy='0'
set openclash.config.custom_host='0'
set openclash.config.global_client_fingerprint='0'
set openclash.config.delay_start='30'
set openclash.config.dashboard_type='Meta'
set openclash.config.enable='1'
set openclash.config.yacd_type='Meta'
set openclash.config.restart='0'
set openclash.config.config_reload='1'
set openclash.@dns_servers[0]=dns_servers
set openclash.@dns_servers[0].group='nameserver'
set openclash.@dns_servers[0].type='udp'
set openclash.@dns_servers[0].enabled='1'
set openclash.@dns_servers[0].ip='127.0.0.1'
set openclash.@dns_servers[0].node_resolve='0'
set openclash.@dns_servers[0].interface='Disable'
set openclash.@dns_servers[0].specific_group='Disable'
set openclash.@dns_servers[0].port='6053'
set openclash.@dns_servers[1]=dns_servers
set openclash.@dns_servers[1].type='udp'
set openclash.@dns_servers[1].node_resolve='0'
set openclash.@dns_servers[1].interface='Disable'
set openclash.@dns_servers[1].specific_group='Disable'
set openclash.@dns_servers[1].group='fallback'
set openclash.@dns_servers[1].ip='127.0.0.1'
set openclash.@dns_servers[1].port='6553'
set openclash.@dns_servers[1].enabled='1'
set openclash.@dns_servers[2]=dns_servers
set openclash.@dns_servers[2].group='nameserver'
set openclash.@dns_servers[2].type='udp'
set openclash.@dns_servers[2].ip='119.28.28.28'
set openclash.@dns_servers[2].enabled='0'
set openclash.@dns_servers[2].node_resolve='0'
set openclash.@dns_servers[2].interface='Disable'
set openclash.@dns_servers[2].specific_group='Disable'
set openclash.@dns_servers[3]=dns_servers
set openclash.@dns_servers[3].group='nameserver'
set openclash.@dns_servers[3].type='udp'
set openclash.@dns_servers[3].ip='223.5.5.5'
set openclash.@dns_servers[3].enabled='0'
set openclash.@dns_servers[3].node_resolve='0'
set openclash.@dns_servers[3].interface='Disable'
set openclash.@dns_servers[3].specific_group='Disable'
set openclash.@dns_servers[4]=dns_servers
set openclash.@dns_servers[4].type='https'
set openclash.@dns_servers[4].ip='doh.pub/dns-query'
set openclash.@dns_servers[4].group='nameserver'
set openclash.@dns_servers[4].enabled='0'
set openclash.@dns_servers[4].node_resolve='0'
set openclash.@dns_servers[4].interface='Disable'
set openclash.@dns_servers[4].specific_group='Disable'
set openclash.@dns_servers[5]=dns_servers
set openclash.@dns_servers[5].type='https'
set openclash.@dns_servers[5].ip='dns.alidns.com/dns-query'
set openclash.@dns_servers[5].group='nameserver'
set openclash.@dns_servers[5].enabled='0'
set openclash.@dns_servers[5].node_resolve='0'
set openclash.@dns_servers[5].interface='Disable'
set openclash.@dns_servers[5].specific_group='Disable'
set openclash.@dns_servers[6]=dns_servers
set openclash.@dns_servers[6].group='fallback'
set openclash.@dns_servers[6].node_resolve='0'
set openclash.@dns_servers[6].interface='Disable'
set openclash.@dns_servers[6].specific_group='Disable'
set openclash.@dns_servers[6].ip='127.0.0.1'
set openclash.@dns_servers[6].port='7553'
set openclash.@dns_servers[6].type='udp'
set openclash.@dns_servers[6].enabled='0'
set openclash.@dns_servers[7]=dns_servers
set openclash.@dns_servers[7].group='fallback'
set openclash.@dns_servers[7].ip='dns.google'
set openclash.@dns_servers[7].port='853'
set openclash.@dns_servers[7].type='tls'
set openclash.@dns_servers[7].enabled='0'
set openclash.@dns_servers[7].node_resolve='0'
set openclash.@dns_servers[7].interface='Disable'
set openclash.@dns_servers[7].specific_group='Disable'
set openclash.@dns_servers[8]=dns_servers
set openclash.@dns_servers[8].group='fallback'
set openclash.@dns_servers[8].type='https'
set openclash.@dns_servers[8].ip='1.1.1.1/dns-query'
set openclash.@dns_servers[8].enabled='0'
set openclash.@dns_servers[8].node_resolve='0'
set openclash.@dns_servers[8].interface='Disable'
set openclash.@dns_servers[8].specific_group='Disable'
set openclash.@dns_servers[9]=dns_servers
set openclash.@dns_servers[9].group='fallback'
set openclash.@dns_servers[9].ip='1.1.1.1'
set openclash.@dns_servers[9].port='853'
set openclash.@dns_servers[9].type='tls'
set openclash.@dns_servers[9].enabled='0'
set openclash.@dns_servers[9].node_resolve='0'
set openclash.@dns_servers[9].interface='Disable'
set openclash.@dns_servers[9].specific_group='Disable'
set openclash.@dns_servers[10]=dns_servers
set openclash.@dns_servers[10].enabled='0'
set openclash.@dns_servers[10].group='fallback'
set openclash.@dns_servers[10].ip='8.8.8.8'
set openclash.@dns_servers[10].port='853'
set openclash.@dns_servers[10].type='tls'
set openclash.@dns_servers[10].node_resolve='0'
set openclash.@dns_servers[10].interface='Disable'
set openclash.@dns_servers[10].specific_group='Disable'
set openclash.@dns_servers[11]=dns_servers
set openclash.@dns_servers[11].type='udp'
set openclash.@dns_servers[11].group='fallback'
set openclash.@dns_servers[11].ip='2001:4860:4860::8888'
set openclash.@dns_servers[11].port='53'
set openclash.@dns_servers[11].enabled='0'
set openclash.@dns_servers[11].node_resolve='0'
set openclash.@dns_servers[11].interface='Disable'
set openclash.@dns_servers[11].specific_group='Disable'
set openclash.@dns_servers[12]=dns_servers
set openclash.@dns_servers[12].type='udp'
set openclash.@dns_servers[12].group='fallback'
set openclash.@dns_servers[12].ip='2001:4860:4860::8844'
set openclash.@dns_servers[12].port='53'
set openclash.@dns_servers[12].enabled='0'
set openclash.@dns_servers[12].node_resolve='0'
set openclash.@dns_servers[12].interface='Disable'
set openclash.@dns_servers[12].specific_group='Disable'
set openclash.@dns_servers[13]=dns_servers
set openclash.@dns_servers[13].type='udp'
set openclash.@dns_servers[13].group='fallback'
set openclash.@dns_servers[13].ip='2001:da8::666'
set openclash.@dns_servers[13].port='53'
set openclash.@dns_servers[13].enabled='0'
set openclash.@dns_servers[13].node_resolve='0'
set openclash.@dns_servers[13].interface='Disable'
set openclash.@dns_servers[13].specific_group='Disable'
set openclash.@dns_servers[14]=dns_servers
set openclash.@dns_servers[14].group='fallback'
set openclash.@dns_servers[14].type='https'
set openclash.@dns_servers[14].ip='public.dns.iij.jp/dns-query'
set openclash.@dns_servers[14].enabled='0'
set openclash.@dns_servers[14].node_resolve='0'
set openclash.@dns_servers[14].interface='Disable'
set openclash.@dns_servers[14].specific_group='Disable'
set openclash.@dns_servers[15]=dns_servers
set openclash.@dns_servers[15].group='fallback'
set openclash.@dns_servers[15].type='https'
set openclash.@dns_servers[15].ip='jp.tiar.app/dns-query'
set openclash.@dns_servers[15].enabled='0'
set openclash.@dns_servers[15].node_resolve='0'
set openclash.@dns_servers[15].interface='Disable'
set openclash.@dns_servers[15].specific_group='Disable'
set openclash.@dns_servers[16]=dns_servers
set openclash.@dns_servers[16].group='fallback'
set openclash.@dns_servers[16].type='https'
set openclash.@dns_servers[16].ip='jp.tiarap.org/dns-query'
set openclash.@dns_servers[16].enabled='0'
set openclash.@dns_servers[16].node_resolve='0'
set openclash.@dns_servers[16].interface='Disable'
set openclash.@dns_servers[16].specific_group='Disable'
set openclash.@dns_servers[17]=dns_servers
set openclash.@dns_servers[17].group='fallback'
set openclash.@dns_servers[17].ip='jp.tiar.app'
set openclash.@dns_servers[17].type='tls'
set openclash.@dns_servers[17].enabled='0'
set openclash.@dns_servers[17].node_resolve='0'
set openclash.@dns_servers[17].interface='Disable'
set openclash.@dns_servers[17].specific_group='Disable'
set openclash.@dns_servers[18]=dns_servers
set openclash.@dns_servers[18].group='fallback'
set openclash.@dns_servers[18].ip='dot.tiar.app'
set openclash.@dns_servers[18].type='tls'
set openclash.@dns_servers[18].enabled='0'
set openclash.@dns_servers[18].node_resolve='0'
set openclash.@dns_servers[18].interface='Disable'
set openclash.@dns_servers[18].specific_group='Disable'
set openclash.@config_subscribe[0]=config_subscribe
set openclash.@config_subscribe[0].enabled='1'
set openclash.@config_subscribe[0].name='DUANG'
set openclash.@config_subscribe[0].emoji='false'
set openclash.@config_subscribe[0].udp='true'
set openclash.@config_subscribe[0].skip_cert_verify='false'
set openclash.@config_subscribe[0].sort='false'
set openclash.@config_subscribe[0].node_type='false'
set openclash.@config_subscribe[0].rule_provider='false'
set openclash.@config_subscribe[0].sub_convert='0'
set openclash.@config_subscribe[0].keyword='香港' '新加坡'
set openclash.@config_subscribe[0].address='$SUBSCRIBE'
set openclash.@authentication[0]=authentication
set openclash.@authentication[0].enabled='1'
set openclash.@authentication[0].username='$username'
set openclash.@authentication[0].password='$password'

commit
EOF